# cron.d/capteur
@daily root /opt/log-file-cleaning 
0 1 * * * root /usr/bin/find /var/log/ -type f -regex '.*\.[0-9]+\.gz$' -delete
0 1 * * * root /usr/bin/find /opt/ -type f -regex '.*\.[0-9]$' -delete
* * * * * root /opt/boot.sh
