# cron.d/capteur
@daily /opt/log-file-cleaning 
0 1 * * *	/usr/bin/find /var/log/ -type f -regex '.*\.[0-9]+\.gz$' -delete
0 1 * * *	/usr/bin/find /opt/ -type f -regex '.*\.[0-9]$' -delete
* * * * * 	/opt/boot.sh
