# cron.d/capteur
@reboot root /opt/capteur/cron.sh
@hourly root /opt/capteur/cron.sh
