# cron.d/capteur
*/1 * * * * root /bin/sh /opt/cron.sh  > /dev/null
